# web/docker/web2py/web2py.conf
[program:web2py]
command=python3 /root/web2py/anyserver.py -s gunicorn -i 0.0.0.0 -p 8000 -w 3
autostart=true
autorestart=true
